# admin_chj/crontab/autooutputpricenews.sh
#!/bin/bash 
/usr/bin/php /data/www/chinacarbon_new/crontab/autoputmarketnews.php